# master/grabdropbox.service
[Unit]
Description=Grab new messages files from dropbox service

[Service]
Type=simple
ExecStart=/usr/local/bin/grabdropbox.sh

